# 031-aws-cloud/01-hello/mhf.service
[Unit]
Description=Go Server

[Service]
ExecStart=/home/ubuntu/mybinary
User=root
Group=root
Restart=always

[Install]
WantedBy=multi-user.target